# crontab
10 13 * * * uv --directory /app run python -m sumodl > /proc/1/fd/1 2>/proc/1/fd/2
